# client.ini
[program:client]
command=python client.py --master-hostname sscloud11 --master-port 5673 --thread-number 150 --decode-path /cephfs2/mr_decode/ --unzip-path /mem_swap/ --client-name `hostname`
directory=/home/enoss/mr_parse/
stopasgroup=true
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/home/enoss/mr_parse/logs/client.log
stderr_logfile=/home/enoss/mr_parse/logs/client.err
