# master.ini
[program:master]
command=bash start_master.sh
directory=/home/enoss/mr_parse/
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=false
stdout_logfile=/home/enoss/mr_parse/logs/master.log
stderr_logfile=/home/enoss/mr_parse/logs/master.err
